# examples/service/empower-ynab-import.timer
[Unit]
Description=Empower YNAB Import Timer

[Timer]
OnCalendar=*-*-* 6:00:00
Persistent=false

[Install]
WantedBy=timers.target